# build/filesystem/etc/systemd/system/post_update_script.service
[Unit]
Description=Run my update script once after reboot
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/etc/post_update/post_update.sh
# Condition to check if the script has already run
ConditionPathExists=/etc/post_update/post_update.flag

[Install]
WantedBy=multi-user.target